# docker/supervisor/chat-gpt.conf
[supervisord]
nodaemon=true

[program:frontend]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:backend]
command=/backend
autostart=true
autorestart=true